# munpanel.cron.example
* * * * * php /var/www/munpanel/artisan schedule:run >> /dev/null 2>&1
